# parakaleo.service
#!/bin/bash
[Unit]
Description=Parakaleo Medical Clinic App
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/parakaleo
27a7290 (WIP: local changes before rebase)
ExecStart=/home/pi/parakaleo/start_parakaleo.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
